# supervisord_logparser.conf
[program:logparser]
command=/usr/local/bin/python /usr/local/bin/logparser -dir /opt/logs
autostart=true
startreties=3
autorestart=unexpected
