# supervisord_scrapyd.conf
[program:scrapyd]
command=/usr/local/bin/python /usr/local/bin/scrapyd
autostart=true
startreties=3
autorestart=unexpected
